# systemd/lcd12864.service
[Unit]
Description=Lcd12864 Daemon

[Service]
ExecStart=/home/xiaolou/workspace/gpio/lcd -d
Type=forking

[Install]
WantedBy=multi-user.target
